# server/supervisord.conf
[supervisord]
nodaemon=true

[program:jupyter]
command=/usr/local/bin/start-notebook.sh --NotebookApp.token=

[program:apache2]
command=/usr/sbin/apachectl -DFOREGROUND
